# root/profile.d/java.sh
# fix ugly font
export _JAVA_OPTIONS='-Dawt.useSystemAAFontSettings=gasp'

