# root/etc/systemd/system/efistub-intelucode.path
[Unit]
Description=Copy intelucode to UEFISYS Partition

[Path]
PathChanged=/boot/intel-ucode.img

[Install]
WantedBy=multi-user.target
WantedBy=system-update.target

